# vibe-work-bot.service
[Unit]
Description=Vibe Work Bot - Telegram Bot for Construction Reporting
After=network.target

[Service]
Type=simple
User=<your-user>
WorkingDirectory=/path/to/vibe_work_bot
ExecStart=/path/to/vibe_work_bot/venv/bin/python /path/to/vibe_work_bot/main.py
Restart=on-failure
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=vibe-work-bot
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
